# install/retypexd.service
[Unit]
Description=retypex keyboard layout switcher daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/retypexd
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
